# container-files/etc/services.d/syslog-ng/run
#!/bin/execlineb -P

/sbin/syslog-ng --cfgfile=/etc/syslog-ng/syslog-ng.conf --no-caps -F
